# systemd/.config/systemd/user/avantree-webrtc-link.service
[Unit]
Description=Avantree WebRTC Processing Auto-Link Service
After=pipewire.service wireplumber.service
BindsTo=pipewire.service
RequisiteOf=pipewire.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/bin/avantree-webrtc-connect
ExecStop=/bin/bash -c 'if [ -f /tmp/avantree-webrtc-module-id ]; then pactl unload-module "$(cat /tmp/avantree-webrtc-module-id)" 2>/dev/null || true; rm -f /tmp/avantree-webrtc-module-id; fi'
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=pipewire.service